# myxrandr2/90-myxrandr2d.rules
# watch for display changes
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/local/myxrandr2d change"

